# supervisord.conf
[supervisord]
logfile=/dev/stdout
nodaemon=true
user=root

[program:app]
command=

[program:scraper]